# WeMo_Link/files/usr/lib/opkg/info/sysevent.list
/etc/init.d/syseventd
/sbin/sysevent
/sbin/syseventd_fork_helper
/sbin/syseventd
/lib/libsysevent.so
