# WeMo_Insight/files/usr/lib/opkg/info/HalInsight.list
/etc/init.d/insightd
/sbin/insightd
/lib/libinsight.so
